# conf/supervisord.conf
[program:wedding-core]
command=python manage.py runserver
directory=/dir/to/app
